# luci-theme-pink/root/etc/uci-defaults/30_luci-theme-pink
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Pink=/luci-static/pink
	set luci.main.mediaurlbase=/luci-static/pink
	commit luci
EOF
exit 0
